# deploy/vps/hushrenewal-agent@.service
[Unit]
Description=HushRenewal %i negotiation agent
After=network-online.target hushrenewal-backend.service
Wants=network-online.target
Requires=hushrenewal-backend.service

[Service]
Type=simple
User=hushrenewal-%i
Group=hushrenewal-%i
WorkingDirectory=/opt/hushrenewal/current/agents
EnvironmentFile=/etc/hushrenewal/agent-%i.defaults
EnvironmentFile=/etc/hushrenewal/groq.env
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/opt/hushrenewal/venvs/agents/bin/uvicorn negotiator.main:app --host 127.0.0.1 --port ${PORT}
Restart=on-failure
RestartSec=5s
TimeoutStopSec=30s
UMask=0077

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
